# templates/mumble/mumble-server-certificate-reloader.timer.j2
[Unit]
Description=Mumble certificate reloader

[Timer]
OnCalendar=*-*-* *:13:37
Persistent=true

[Install]
WantedBy=timers.target
